# fwDCS/deploy/dcs.service
[Unit]
Description=Faberwork DCS (datacenter collector/forwarder)
# DCS needs the network and its TimescaleDB. If the DB isn't ready yet, DCS
# exits at boot (ApplyRetention is fatal) and systemd restarts it (Restart=always)
# until the DB accepts connections — so a hard ordering dep is not required.
# If Postgres runs natively, add it explicitly:  After=postgresql.service
# If Postgres runs in Docker, add:                After=docker.service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# __INSTALL_DIR__ = the fwDCS dir on the VM (holds dcs_prod.sh, build/, dcs.prod.yaml).
# Reuse the existing prod launcher so binary-pick + GOMEMLIMIT + config stay in one
# place; don't duplicate them here.
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/dcs_prod.sh
User=root

# ── Auto-restart on crash OR clean exit ──────────────────────────────────────
Restart=always
RestartSec=3
# Don't hammer-loop forever if wedged: max 5 restarts / 60s, then stop and
# surface it (check: systemctl status dcs / journalctl -u dcs).
StartLimitIntervalSec=60
StartLimitBurst=5

# ── Resource guards ──────────────────────────────────────────────────────────
# Hard cgroup ceiling so a forwarder burst can never take the whole VM down: the
# kernel caps THIS service instead of OOM-killing a random VM process. Set ABOVE
# GOMEMLIMIT (384MiB in dcs_prod.sh) so the Go runtime GCs hard and reclaims
# before the cgroup limit is hit; MemoryHigh applies reclaim pressure first, then
# MemoryMax is the hard stop. On a true OOM systemd restarts (Restart=always).
MemoryHigh=448M
MemoryMax=512M
# Cap CPU so a GC-thrash spike can't starve the rest of a small (2-core) VM.
CPUQuota=80%
Nice=10

[Install]
WantedBy=multi-user.target
