# fwDCS/deploy/dcs.service
[Unit]
Description=Faberwork DCS (datacenter collector/forwarder)
# DCS needs the network and its TimescaleDB. If the DB isn't ready yet, DCS
# exits at boot (ApplyRetention is fatal) and systemd restarts it (Restart=always)
# until the DB accepts connections — so a hard ordering dep is not required.
# If Postgres runs natively, add it explicitly:  After=postgresql.service
# If Postgres runs in Docker, add:                After=docker.service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# __INSTALL_DIR__ = the fwDCS dir on the VM (holds dcs_prod.sh, build/, dcs.prod.yaml).
# Reuse the existing prod launcher so binary-pick + GOMEMLIMIT + config stay in one
# place; don't duplicate them here.
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/dcs_prod.sh
User=root

# ── Auto-restart on crash OR clean exit ──────────────────────────────────────
Restart=always
RestartSec=3
# Don't hammer-loop forever if wedged: max 5 restarts / 60s, then stop and
# surface it (check: systemctl status dcs / journalctl -u dcs).
StartLimitIntervalSec=60
StartLimitBurst=5

# ── Resource guards ──────────────────────────────────────────────────────────
# Hard cgroup ceiling so a runaway can never take the whole VM down: the kernel
# caps THIS service instead of OOM-killing a random VM process. MUST sit WELL
# above GOMEMLIMIT (384MiB in dcs_prod.sh) — GOMEMLIMIT bounds the Go HEAP, but
# RSS = heap + goroutine stacks + runtime/off-heap overhead (~100-150MB more), so
# a cap near GOMEMLIMIT (the old 512M) OOM-kills DCS under normal load → restart
# loop → migrations stall + ingest stops. ~2x GOMEMLIMIT leaves ample headroom
# for normal operation while still catching a genuine runaway. GOMEMLIMIT remains
# the real governor; this is only a backstop. On a true OOM systemd restarts.
MemoryHigh=640M
MemoryMax=768M
# Cap CPU so a GC-thrash spike can't starve the rest of a small (2-core) VM.
CPUQuota=80%
Nice=10

[Install]
WantedBy=multi-user.target
